# services/telemuxd/run
#!/data/data/com.termux/files/usr/bin/sh

SCRIPT_DIR=$(dirname "$(realpath "$0")")

pm2 start "$SCRIPT_DIR/plugin/run.js" --name "telemux"

pm2 save
