# services/telemuxd/run
#!/data/data/com.termux/files/usr/bin/sh

exec node ./plugin/run.js