# services/telemuxd/run
#!/data/data/com.termux/files/usr/bin/sh

exec ts-node ./plugin/run.ts