# services/telemuxd/run
#!/data/data/com.termux/files/usr/bin/sh

SCRIPT_DIR=$(dirname "$(realpath "$0")")

exec node "$SCRIPT_DIR/plugin/index.js"
